# runit/docstore-sharelatex/run
#!/bin/bash
export SHARELATEX_CONFIG=/etc/sharelatex/settings.js
export SESSION_SECRET="${SESSION_SECRET:-change-this-session-secret-in-production}"

NODE_PARAMS=""
if [ "$DEBUG_NODE" == "true" ]; then
    echo "running debug - docstore"
    NODE_PARAMS="--inspect=0.0.0.0:30160"
fi

exec /sbin/setuser www-data env NODE_OPTIONS="$NODE_PARAMS" /usr/bin/yarn --cwd /var/www/sharelatex workspace @overleaf/docstore start >> /var/log/sharelatex/docstore.log 2>&1
